# systemd/codexapp-history-sync.service
[Unit]
Description=CodexApp history sync
After=network-online.target codexapp-web.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=CODEX_HISTORY_VSCODE_HOME=/var/lib/codexapp-linux-web/import-source/.codex
Environment=CODEX_HISTORY_CODEXAPP_HOME=/var/lib/codexapp-linux-web/.codex
Environment=CODEX_HISTORY_SYNC_INTERVAL=30
Environment=CODEX_HISTORY_SYNC_MIN_FREE_MB=10240
ExecStart=/opt/codexapp-linux-web/sync-history.sh
Restart=always
RestartSec=2
StandardOutput=append:/var/log/codexapp-linux-web/history-sync/service.log
StandardError=append:/var/log/codexapp-linux-web/history-sync/service.log

[Install]
WantedBy=multi-user.target
