# systemd/codexapp-app-server.service
[Unit]
Description=CodexApp app-server (persistent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/codexapp-linux-web
Environment=HOME=/var/lib/codexapp-linux-web
Environment=USER=root
Environment=NODE_ENV=production
Environment=CODEX_HOME=/var/lib/codexapp-linux-web/.codex
Environment=CODEXAPP_CODEX_CLI=/usr/bin/codex
Environment=CODEXAPP_APP_VERSION=0.141.0
Environment=CODEXAPP_APP_SERVER_LISTEN=ws://127.0.0.1:12911
Environment=CODEXAPP_ADOPT_EXISTING_APP_SERVER=1
EnvironmentFile=-/etc/codexapp-linux-web/codexapp.env
ExecStart=/opt/codexapp-linux-web/scripts/codexapp-app-server-supervisor.sh
Restart=always
RestartSec=3
KillMode=process
TimeoutStopSec=30s
StandardOutput=append:/var/log/codexapp-linux-web/app-server/service.log
StandardError=append:/var/log/codexapp-linux-web/app-server/service.log
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
